# bookstack/rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Bookstack
# Runs the MySQL database server
# ==============================================================================
if bashio::config.equals 'database_location' 'remote';then
sleep infinity
fi

declare -a options

options+=(--"datadir=/data/mysql")
options+=(--"tmpdir=/tmp/")
options+=(--"user=mysql")
options+=(--"log-error=/data/mysql/log.err")
bashio::log.info "Starting MySQL..."
exec /usr/bin/mysqld_safe "${options[@]}"
