# bookstack/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Bookstack
# Runs the PHP-FPM daemon
# ==============================================================================
declare key
# Wait for MySQL to become available
for _ in $(seq 1 30); do
    if echo 'SELECT 1' | mysql &> /dev/null; then
        break
    fi
    sleep 1
done

key=$(cat /data/bookstack/appkey.txt)
export APP_KEY=${key}
export DB_HOST=localhost
export DB_DATABASE=bookstack
export DB_USERNAME=bookstack
export DB_PASSWORD=bookstack

bashio::log.info "Installing/updating Database"
php /var/www/bookstack/artisan migrate --force


bashio::log.info "Starting PHP-FPM..."

exec php-fpm7 --nodaemonize
